# medium/base-supervisord-httpd-php-mysqld/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/bin/bash -c "/docker-entrypoint.sh"

[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
