# hard/deneb/files/supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command=/usr/sbin/crond -f

[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
